# app/crontab
*/5 * * * * python /app/task.py >> /app/logfile.txt